# .archive/.config/autostart/mpd-mpris.desktop
[Desktop Entry]
Name=mpd-mpris
Comment=MPD mpris daemon
Exec=/run/current-system/sw/bin/mpd-mpris
Type=Application
Terminal=false
NoDisplay=true
